# bin/les-projets-cagnottes-slack.service
[Unit]
Description=Les Projets Cagnottes - Slack Events Catcher Daemon

[Service]
Type=forking
Environment=HOME=/root
ExecStart=/opt/les-projets-cagnottes/slack-events-catcher/current/bin/les-projets-cagnottes-slack.sh start
ExecStop=/opt/les-projets-cagnottes/slack-events-catcher/current/bin/les-projets-cagnottes-slack.sh stop
ExecReload=/opt/les-projets-cagnottes/slack-events-catcher/current/bin/les-projets-cagnottes-slack.sh reload

[Install]
WantedBy=multi-user.target