# bin/slack-events-catcher.service
[Unit]
Description=slack-events-catcher daemon

[Service]
Type=forking
Environment=HOME=/root
ExecStart=/opt/les-projets-cagnottes/slack-events-catcher/bin/slack-events-catcher.sh start
ExecStop=/opt/les-projets-cagnottes/slack-events-catcher/bin/slack-events-catcher.sh stop
ExecReload=/opt/les-projets-cagnottes/slack-events-catcher/bin/slack-events-catcher.sh reload

[Install]
WantedBy=multi-user.target